# install/s6-overlay/s6-rc.d/juicity-server/run
#!/bin/sh
exec /usr/local/bin/juicity-server run -c /usr/local/etc/juicity/server.json
