# install/s6-overlay/s6-rc.d/juicity-client/run
#!/bin/sh
exec /usr/local/bin/juicity-client run -c /usr/local/etc/juicity/client.json
